# deploy/connect5-ws.service
[Unit]
Description=connect5 WebSocket game server
After=network.target

[Service]
Type=simple
User=CONNECT5_USER_PLACEHOLDER
Group=CONNECT5_USER_PLACEHOLDER
WorkingDirectory=CONNECT5_ROOT_PLACEHOLDER
Environment=NODE_ENV=production
Environment=PORT=3003
ExecStart=/usr/bin/npm run server
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
